# containers/make-container@.service
[Unit]
Description=Prepare for Container %i

[Service]
Type=oneshot
ExecStart=/usr/bin/make-container %i

[Install]
WantedBy=multi-user.target

