# router/bind/dyndns.service
[Unit]
Description=Update DNS

[Service]
Type=oneshot
ExecStart=/usr/bin/dyndns
