# router/bind/dyndns.timer
[Unit]
Description=Run nsupdate on a schedule

[Timer]
OnCaledar=minutely
Persistent=true

[Install]
WantedBy=timers.target
