# web/ffserver.service
[Service]
ExecStart=/usr/bin/ffserver
Restart=always

[Install]
WantedBy=multi-user.target

